# deploy/linux-arm64/waddle-view.service
[Unit]
Description=Waddle View TV dashboard
After=graphical-session.target

[Service]
Type=simple
Environment=WADDLE_API_KEY_FILE=/etc/waddle-view/api.key
Environment=DISPLAY=:0
WorkingDirectory=/opt/waddle-view/bundle
ExecStart=/opt/waddle-view/bundle/waddle_display
Restart=on-failure

[Install]
WantedBy=default.target
